# backend/nation-radar-api.service
[Unit]
Description=Nation Radar VPS API Server
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/nation-radar/backend
Environment=VPS_API_PORT=5001
ExecStart=/root/nation-radar/backend/start_api_server.sh
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
