# backend/nation-radar-pipeline-weekly.service
[Unit]
Description=Nation Radar Data Pipeline - Weekly Collection
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/nation-radar/backend
Environment=PATH=/root/nation-radar/venv/bin
ExecStart=/root/nation-radar/backend/run_pipeline.py
Restart=always
RestartSec=604800
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
